# calamaris.crontab
PATH = /sbin:/bin:/usr/sbin:/usr/bin

30 0 * * *	root	cat /var/log/squid/access.log.0 | \
                	nice -39 calamaris -amH 'Squid log - daily' | \
                	mail root
0 3 * * 7	root	cat /var/log/squid/access.log.[0-6] | \
                	nice -39 calamaris -amH 'Squid log - weekly' | \
                	mail root
