# calamaris.crontab
PATH=/sbin:/bin:/usr/sbin:/usr/bin

00 4 * * *	root	calamaris-croniface daily
30 4 * * 0	root	calamaris-croniface weekly
